# ldap2dns.xml
<?xml version="1.0"?>
<!DOCTYPE service_bundle SYSTEM "/usr/share/lib/xml/dtd/service_bundle.dtd.1">
<!--
	Copyright 2006-2011 Alkaloid Networks, LLC
        Licensed under the terms of the GPL version 2. See
        http://www.gnu.org/licenses/old-licenses/gpl-2.0.html

This program is free software; you can redistribute it and/or
modify it under the terms of the GNU General Public License
as published by the Free Software Foundation; either version 2
of the License, or (at your option) any later version.

This program is distributed in the hope that it will be useful,
but WITHOUT ANY WARRANTY; without even the implied warranty of
MERCHANTABILITY or FITNESS FOR A PARTICULAR PURPOSE.  See the
GNU General Public License for more details.

You should have received a copy of the GNU General Public License
along with this program; if not, write to the Free Software
Foundation, Inc., 51 Franklin Street, Fifth Floor, Boston, MA  02110-1301, USA.

	$Id$
-->
<service_bundle type='manifest' name='ALKAldap2dns:ldap2dns'>

<service
	name='application/management/ldap2dns'
	type='service'
	version='1'>

	<exec_method
	    type='method'
	    name='stop'
	    exec=':kill'
	    timeout_seconds='60' />

	<instance name='default' enabled='false' >
		<exec_method
		    type='method'
		    name='start'
		    exec='/usr/bin/ldap2dnsd'
		    timeout_seconds='60' >
			<method_context>
				<method_environment>
					<envvar name='LDAP2DNS_UPDATE' value='60' />
					<envvar name='LDAP2DNS_OUTPUT' value='data' />
				</method_environment>
			</method_context>
		</exec_method>
	</instance>

	<stability value='Unstable' />

	<template>
		<common_name>
			<loctext xml:lang='C'>
				LDAP to DNS manager
			</loctext>
		</common_name>
		<documentation>
			<manpage title='ldap2dns' section='1'
			    manpath='/usr/man' />
		</documentation>
	</template>

</service>

</service_bundle>
